# nanoscreen.service
[Unit]
Description=Adafruit PiTFT NanoScreen Service

[Service]
Type=simple
ExecStart=/usr/local/bin/nanoscreen

[Install]
WantedBy = basic.target
