# square.service
[Unit]
Description=Adafruit PiTFT Square Service

[Service]
Type=simple
ExecStart=/usr/local/bin/square

[Install]
WantedBy = basic.target
