# tftcp.service
[Unit]
Description=Adafruit PiTFT TFTCP Service

[Service]
Type=simple
ExecStart=/usr/local/bin/tftcp

[Install]
WantedBy = basic.target
